# distrib/systemd/goxlr-initialiser.service
[Unit]
Description=GoXLR Initialiser

[Service]
Type=oneshot
ExecStart=/usr/bin/goxlr-initialiser

[Install]
WantedBy=graphical.target
